# .crontab
* * * * * printf '# as of ' > ~/.crontab ; date >> ~/.crontab ; crontab -l >> ~/.crontab
