# bash/bin/profile.d/umask.sh
#!/bin/sh

# -rw-rw-r--, drwxrwxr-x
umask 0002
