# bin/profile.d/disable-apple-terminal-session-history.sh
#!/bin/bash
touch ~/.bash_sessions_disable
export SHELL_SESSION_HISTORY=0
